# scripts/invo-esp32-fota-server.service
[Unit]
Description=INVO ESP32 FOTA HTTP Server
Documentation=Serves esp32-fota/ dir on :8080 for plain-HTTP OTA (avoids SDIO/TLS issues)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=intelli
WorkingDirectory=/home/intelli/esp32-fota
ExecStart=/usr/bin/python3 -m http.server 8080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
